# crontab
SHELL=/bin/bash
*/5 * * * * root . /var/www/zohoint/env.sh; /usr/local/bin/python /var/www/zohoint/main_zoho2jira.py > /var/log/cron.log
