# crontab
SHELL=/bin/bash
*/5 * * * * root /usr/local/bin/python /var/www/zohoint/main_zoho2jira.py > /var/log/cron.log 2>&1

